# vm/var/lib/dpkg/info/courier-imap.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/courier-imap
/usr/share/doc/courier-imap/changelog.gz
/usr/share/doc/courier-imap/README.proxy.html
/usr/share/doc/courier-imap/README.imapkeywords.html
/usr/share/doc/courier-imap/README.Debian
/usr/share/doc/courier-imap/copyright
/usr/share/doc/courier-imap/NEWS.Debian.gz
/usr/share/doc/courier-imap/README.proxy.gz
/usr/share/doc/courier-imap/changelog.Debian.gz
/usr/share/courier
/usr/share/courier/webadmin
/usr/share/courier/webadmin/admin-40imap.html
/usr/share/courier/webadmin/admin-40imap.pl
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/courier-imap
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/imapd.8.gz
/usr/bin
/usr/bin/imapd
/usr/lib
/usr/lib/courier
/usr/lib/courier/imapd
/usr/lib/courier/courier
/usr/lib/courier/courier/imaplogin
/usr/sbin
/etc
/etc/courier
/etc/courier/imapd
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/courier-imap
/etc/logcheck/violations.ignore.d
/etc/logcheck/violations.ignore.d/courier-imap
/etc/pam.d
/etc/pam.d/imap
/etc/init.d
/etc/init.d/courier-imap
/usr/sbin/imapd
