# vm/var/lib/dpkg/info/klogd.list
/etc
/etc/init.d
/etc/init.d/klogd
/etc/default
/etc/default/klogd
