# vm/var/lib/dpkg/info/sysklogd.list
/etc
/etc/init.d
/etc/init.d/sysklogd
/etc/cron.daily
/etc/cron.daily/sysklogd
/etc/cron.weekly
/etc/cron.weekly/sysklogd
/etc/default
/etc/default/syslogd
/etc/syslog.conf
