# vm/var/lib/dpkg/info/tomcat6.list
/.
/var
/var/log
/var/log/tomcat6
/var/lib
/var/lib/tomcat6
/var/lib/tomcat6/webapps
/var/lib/tomcat6/common
/var/lib/tomcat6/common/classes
/var/lib/tomcat6/server
/var/lib/tomcat6/server/classes
/var/lib/tomcat6/shared
/var/lib/tomcat6/shared/classes
/var/cache
/var/cache/tomcat6
/etc
/etc/tomcat6
/etc/tomcat6/Catalina
/etc/tomcat6/Catalina/localhost
/etc/tomcat6/catalina.properties
/etc/tomcat6/logging.properties
/etc/tomcat6/context.xml
/etc/tomcat6/server.xml
/etc/tomcat6/tomcat-users.xml
/etc/tomcat6/web.xml
/etc/tomcat6/policy.d
/etc/tomcat6/policy.d/01system.policy
/etc/tomcat6/policy.d/02debian.policy
/etc/tomcat6/policy.d/03catalina.policy
/etc/tomcat6/policy.d/04webapps.policy
/etc/tomcat6/policy.d/50local.policy
/etc/default
/etc/default/tomcat6
/etc/init.d
/etc/init.d/tomcat6
/etc/cron.daily
/etc/cron.daily/tomcat6
/usr
/usr/share
/usr/share/doc
/usr/share/doc/tomcat6
/usr/share/doc/tomcat6/copyright
/usr/share/doc/tomcat6/changelog.Debian.gz
/usr/share/tomcat6
/usr/share/tomcat6/webapps
/usr/share/tomcat6/webapps/default_root
/usr/share/tomcat6/webapps/default_root/META-INF
/usr/share/tomcat6/webapps/default_root/META-INF/context.xml
/usr/share/tomcat6/webapps/default_root/index.html
/var/lib/tomcat6/logs
/var/lib/tomcat6/work
/var/lib/tomcat6/conf
/usr/share/doc/tomcat6/README.Debian.gz
